# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#
# https://github.com/hosac | hosac@gmx.net
#
# Reload ALSA Audio Mixer settings, otherwise no sound...
alsactl restore -f /etc/asound.state

exit 0
